# src/agentic/com.zoom.daily-pipeline.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.zoom.daily-pipeline</string>

    <key>Program</key>
    <string>/bin/bash</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/abhinavs./Documents/Zoom/src/agentic/daily_pipeline.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/abhinavs./Documents/Zoom</string>

    <key>RunAtLoad</key>
    <false/>

    <!-- Daily at 18:00 IST = 12:30 UTC; weekday only (Mon=1..Fri=5) -->
    <!-- Repeated for each weekday so it doesn't fire on weekends -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>18</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>18</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>18</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>18</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>18</integer><key>Minute</key><integer>0</integer></dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/abhinavs./Documents/Zoom/logs/launchd_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/abhinavs./Documents/Zoom/logs/launchd_stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/bin:/bin:/usr/sbin:/sbin:/usr/local/bin:/opt/homebrew/bin</string>
        <key>PYTHONPATH</key>
        <string>/Users/abhinavs./Documents/Zoom</string>
    </dict>
</dict>
</plist>
